# packaging/vexyl-guard.service
[Unit]
Description=Vexyl Guard host agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/vexyl/guard.conf
ExecStart=/usr/local/sbin/vexyl-guard daemon
Restart=always
RestartSec=5s
StateDirectory=vexyl
ConfigurationDirectory=vexyl
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
